# systemd/thornhill-ci-deploy.timer
[Unit]
Description=Poll GitHub for the latest CI-proven Thornhill main revision

[Timer]
OnBootSec=5min
OnUnitInactiveSec=15min
RandomizedDelaySec=1min
Persistent=true
Unit=thornhill-ci-deploy.service

[Install]
WantedBy=timers.target
